# scripts/fix-dns.service.example
# systemd-resolved stub (127.0.0.53) often breaks outbound DNS in incus/LXC containers.
# Install: bash scripts/install-fix-dns.sh
#
[Unit]
Description=Fix container DNS (public resolvers)
Before=network-online.target cursor-voice.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'printf "nameserver 8.8.8.8\nnameserver 1.1.1.1\n" > /etc/resolv.conf'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
